# conf.d/main.conf
[program:main]
command=/app/pwd-recurse/main.py
autorestart=true
